# deploy/systemd/tradifi-dashboard-update.timer
[Unit]
Description=Run TradFi Funding Dashboard update every 8 hours

[Timer]
OnCalendar=*-*-* 00,08,16:17:00 Asia/Shanghai
Persistent=true
AccuracySec=1min
RandomizedDelaySec=0
Unit=tradifi-dashboard-update.service

[Install]
WantedBy=timers.target
